# base-vm/dependencies/skel/.config/autostart/keyboard-language.desktop
[Desktop Entry]
Type=Application
Name=SetKeyboardLanguage
Exec=setxkbmap us
