# base-vm/dependencies/skel/.config/autostart/set-keyboard-langauge.desktop
[Desktop Entry]
Type=Application
Name=SetKeyboardLanguage
Exec=setxkbmap de,us,gb,fr,it,es -option grp:alt_shift_toggle
