# base-vm/dependencies/skel/.config/autostart/show-welcome.desktop
[Desktop Entry]
Type=Application
Name=Welcome-Message
Exec=/usr/share/kx.as.code/showWelcome.sh
